# deployment/systemd/graphqna-slackbot.service
[Unit]
Description=GraphQnA Slack Bot
After=network.target graphqna-api.service
Wants=network-online.target
Requires=graphqna-api.service

[Service]
Type=simple
User=graphqna
WorkingDirectory=/opt/graphqna
EnvironmentFile=/opt/graphqna/.env
ExecStart=/opt/graphqna/venv/bin/python /opt/graphqna/scripts/run_slack_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=graphqna-slackbot

[Install]
WantedBy=multi-user.target