# deployment/systemd/graphqna-api.service
[Unit]
Description=GraphQnA API Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/GraphQnA
EnvironmentFile=/home/ubuntu/GraphQnA/.env
ExecStart=/home/ubuntu/GraphQnA/venv/bin/python /home/ubuntu/GraphQnA/scripts/run_api.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=graphqna-api

[Install]
WantedBy=multi-user.target
