# deployment/systemd/graphqna-slackbot.service
[Unit]
Description=GraphQnA Slack Bot
After=network.target graphqna-api.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/GraphQnA
EnvironmentFile=/home/ubuntu/GraphQnA/.env
ExecStart=/home/ubuntu/GraphQnA/venv/bin/python /home/ubuntu/GraphQnA/scripts/run_slack_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=graphqna-slackbot

[Install]
WantedBy=multi-user.target
